# gaffer-docker/centos6/base/ssh_service_supervisord.conf
[supervisord]
nodaemon=false


[program:sshd]
command=/usr/sbin/sshd -D
